# roles/galera_healthcheck/templates/mysqlchk.socket
[Unit]
Description=Mariadb Galera Cluster node check socket

[Socket]
ListenStream=9200
Accept=yes

[Install]
WantedBy=sockets.target
