# roles/galera_healthcheck/templates/mysqlchk@.service
[Unit]
Description=Galera Cluster node check service

[Service]
ExecStart=/usr/bin/clustercheck
StandardInput=socket
